# release_files/relay/openzro-relay.service
[Unit]
Description=openZro relay server (TURN-style fallback for peers behind symmetric NAT)
Documentation=https://docs.openzro.io/selfhosted
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openzro
Group=openzro

EnvironmentFile=-/etc/default/openzro-relay

# OZ_LISTEN_ADDRESS, OZ_EXPOSED_ADDRESS, OZ_AUTH_SECRET come from
# the env file (see openzro-relay.env). All required.
ExecStart=/usr/bin/openzro-relay

Restart=on-failure
RestartSec=15s

# Relay binds a single TCP port (default :33080) — no privileged
# operations needed once we drop CAP_NET_BIND_SERVICE.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openzro-relay

[Install]
WantedBy=multi-user.target
